# deploy/systemd/gscale-mobileapi.service
[Unit]
Description=GScale Zebra mobile API
After=network-online.target gscale-scale.service
Wants=network-online.target gscale-scale.service

[Service]
Type=simple
User=__APP_USER__
Group=__APP_GROUP__
WorkingDirectory=__PREFIX__
EnvironmentFile=-__PREFIX__/config/mobileapi.env
Environment=MOBILE_API_SETUP_FILE=__PREFIX__/config/core.env
ExecStart=__PREFIX__/bin/mobileapi
Restart=always
RestartSec=1
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
